# dotcloud/supervisord.conf
[program:dcdaemon]
command = python dcdaemon.py
directory = /home/dotcloud/current/
stdout_logfile = /var/log/supervisor/dcdaemon.log
